# supervisord.conf
[supervisord]
nodaemon=true

[program:apache]
command=httpd-foreground
autostart=true
autorestart=true